# deploy/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/run/supervisord.pid
childlogdir=/data/logs

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

__POSTGRES_PROGRAM__
[program:sshd]
command=/usr/sbin/sshd -D -e -f /app/deploy/ssh/sshd_config
directory=/app
autostart=true
autorestart=true
startsecs=1
stdout_logfile=/data/logs/sshd.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true
priority=20

[program:portal]
command=/app/deploy/bin/start-portal.sh
directory=/app
user=pilotswarm
autostart=true
autorestart=true
startsecs=1
stdout_logfile=/data/logs/portal.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true
priority=30

[program:worker-a]
command=/app/deploy/bin/start-worker.sh worker-a
directory=/app
user=pilotswarm
autostart=true
autorestart=true
startsecs=1
stdout_logfile=/data/logs/worker-a.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true
priority=40

[program:worker-b]
command=/app/deploy/bin/start-worker.sh worker-b
directory=/app
user=pilotswarm
autostart=true
autorestart=true
startsecs=1
stdout_logfile=/data/logs/worker-b.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true
priority=41
